# ordx-data-maintain/ord.ini
[program:bitcoind-25.0-main]
priority=1
command=bitcoind -rpcworkqueue=512 -rpcthreads=128 -chain=main -conf=/data/bitcoinData/bitcoin.conf -datadir=/data/bitcoinData/main
autostart=false
autorestart=true

[program:bitcoind-25.0-testnet]
priority=1
command=bitcoind -rpcworkqueue=512 -rpcthreads=128 -chain=test -conf=/data/bitcoinData/bitcoin.conf -datadir=/data/bitcoinData
autostart=false
autorestart=true

[program:ord-0.16.0-main]
priority=2
command=ord0.16.0 --chain mainnet --bitcoin-rpc-url 192.168.1.102:8332 --data-dir /data/ordData/0.16.0/main --bitcoin-rpc-username jacky --bitcoin-rpc-password _RZekaGRgKQJSIOYi6vq0_CkJtjoCootamy81J2cDn0 --first-inscription-height 767430 --index-sats --minify server --http --http-port 82
autostart=false
autorestart=true
user=root
stderr_logfile=/var/log/supervisor/ord0.16.0-mainnet.out.log
stdout_logfile=/var/log/supervisor/ord0.16.0-mainnet.out.log

[program:ord0.16.0-testnet]
priority=2
command=ord0.16.0 --chain testnet --bitcoin-rpc-url 192.168.1.102:18332 --data-dir /data/ordData/0.16.0 --bitcoin-rpc-username jacky --bitcoin-rpc-password 123456 --first-inscription-height 2413343 --index-sats --minify server --http --http-port 80
autostart=false
autorestart=true
user=root
stderr_logfile=/var/log/supervisor/ord0.16.0-testnet.out.log
stdout_logfile=/var/log/supervisor/ord0.16.0-testnet.out.log

[program:ordx-server-master-main]
priority=3
command=ordx-server -env /etc/ordx/test-main.env
autostart=false
autorestart=true
user=root
stderr_logfile=/var/log/supervisor/ordx-server-main.out.log
stdout_logfile=/var/log/supervisor/ordx-server-main.out.log

[program:ordx-server-master-testnet]
priority=3
command=ordx-server -env /etc/ordx/test-testnet.env
autostart=false
autorestart=true
user=root
stderr_logfile=/var/log/supervisor/ordx-server-testnet.out.log
stdout_logfile=/var/log/supervisor/ordx-server-testnet.out.log

[eventlistener:ordx_handler]
command=/root/ordx-data-maintain/supervisor-eventlistener-ordx-handler.py
events=PROCESS_STATE_EXITED
process_name=%(program_name)s_%(process_num)02d
numprocs=5
buffer_size=10
directory=/tmp
umask=022
priority=-1
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
exitcodes=0
stopsignal=QUIT
stopwaitsecs=10
stopasgroup=false
killasgroup=false
redirect_stderr=false
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_events_enabled=false
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stderr_events_enabled=false
serverurl=AUTO
stdout_logfile=/var/log/supervisor/supervisor-eventlistener-ordx-handler.log
stderr_logfile=/var/log/supervisor/supervisor-eventlistener-ordx-handler.log
