# help/supervisor/conf/101.ordx.ini
[program:bitcoind-28.0-main]
priority=1
command=bitcoind -rpcworkqueue=512 -rpcthreads=128 -chain=main -conf=/data/bitcoin-data/bitcoin.conf -datadir=/data/bitcoin-data/main
autostart=true
autorestart=true
topsignal=INT

[program:bitcoind-28.0-testnet4]
priority=1
command=bitcoind -rpcworkqueue=512 -rpcthreads=128 -chain=testnet4 -conf=/data/bitcoin-data/bitcoin.conf -datadir=/data/bitcoin-data
autostart=true
autorestart=true
stopsignal=INT

[program:ordx-mainnet-master]
priority=3
command=ordx-mainnet-master -env /etc/ordx/prd-mainnet-master.env
autostart=true
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ordx-mainnet-master.out.log
stdout_logfile=/var/log/supervisor/ordx-mainnet-master.out.log

[program:ordx-testnet4-master]
priority=3
command=ordx-testnet4-master -env /etc/ordx/prd-testnet4-master.env
autostart=true
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ordx-testnet4-master.out.log
stdout_logfile=/var/log/supervisor/ordx-testnet4-master.out.log

[program:ordx-mainnet-slave]
priority=3
command=ordx-mainnet-slave -env /etc/ordx/prd-mainnet-slave.env
autostart=false
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ordx-mainnet-slave.out.log
stdout_logfile=/var/log/supervisor/ordx-mainnet-slave.out.log

[program:ordx-testnet4-slave]
priority=3
command=ordx-testnet4-slave -env /etc/ordx/prd-testnet4-slave.env
autostart=true
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ordx-testnet4-slave.out.log
stdout_logfile=/var/log/supervisor/ordx-testnet4-slave.out.log

[program:ord-0.19.1-testnet4-testnet4]
priority=2
command=ord0.19.1-testnet4 --chain testnet4 --bitcoin-rpc-url 192.168.10.101:28332 --data-dir /data/ord-data/0.19.1-testnet4 --bitcoin-rpc-username jacky --bitcoin-rpc-password 123456 --first-inscription-height 27228 --index-runes --index-sats server --http --http-port 83
autostart=true
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ord-0.19.1-testnet4-testnet4.out.log
stdout_logfile=/var/log/supervisor/ord-0.19.1-testnet4-testnet4.out.log

[program:ord-0.19.1-testnet4-mainnet]
priority=2
command=ord0.19.1-testnet4 --chain mainnet --bitcoin-rpc-url 192.168.10.101:8332 --data-dir /data/ord-data/0.19.1-testnet4/mainnet --bitcoin-rpc-username jacky --bitcoin-rpc-password _RZekaGRgKQJSIOYi6vq0_CkJtjoCootamy81J2cDn0 --first-inscription-height 767430 --index-runes --index-sats server --http --http-port 84
autostart=true
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ord-0.19.1-testnet4-mainnet.out.log
stdout_logfile=/var/log/supervisor/ord-0.19.1-testnet4-mainnet.out.log