# help/supervisor/conf/102.ordx.ini
[program:bitcoind-28.0-mainnet]
priority=1
command=bitcoind -rpcworkqueue=512 -rpcthreads=128 -chain=main -conf=/data/bitcoin-data/bitcoin.conf -datadir=/data/bitcoin-data/main
autostart=true
autorestart=true
topsignal=INT

[program:bitcoind-28.0-testnet4]
priority=1
command=bitcoind -rpcworkqueue=512 -rpcthreads=128 -chain=testnet4 -conf=/data/bitcoin-data/bitcoin.conf -datadir=/data/bitcoin-data
autostart=true
autorestart=true
stopsignal=INT

[program:ord-0.19.1-mainnet]
priority=2
command=ord0.19.1-testnet4 --chain mainnet --bitcoin-rpc-url 192.168.10.102:8332 --data-dir /data/ord-data/0.19.1 --bitcoin-rpc-username jacky --bitcoin-rpc-password _RZekaGRgKQJSIOYi6vq0_CkJtjoCootamy81J2cDn0 --first-inscription-height 767430 --index-runes --index-sats server --http --http-port 82
autostart=true
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ord-0.19.1-testnet4-mainnet.out.log
stdout_logfile=/var/log/supervisor/ord-0.19.1-testnet4-mainnet.out.log

[program:ord-0.19.1-testnet4]
priority=2
command=ord0.19.1-testnet4 --chain testnet4 --bitcoin-rpc-url 192.168.10.102:28332 --data-dir /data/ord-data/0.19.1 --bitcoin-rpc-username jacky --bitcoin-rpc-password 123456 --first-inscription-height 27228 --index-runes --index-sats server --http --http-port 81
autostart=true
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ord-0.19.1-testnet4-testnet4.out.log
stdout_logfile=/var/log/supervisor/ord-0.19.1-testnet4-testnet4.out.log

[program:ordx-mainnet]
priority=3
directory=/data/github/ordx
command=/data/github/ordx/ordx-mainnet -env /data/github/ordx/.env
autostart=false
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ordx-mainnet.out.log
stdout_logfile=/var/log/supervisor/ordx-mainnet.out.log

[program:ordx-testnet4]
priority=3
directory=/data/github/ordx-testnet
command=/data/github/ordx-testnet/ordx-testnet -env /data/github/ordx-testnet/.env
autostart=false
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ordx-testnet4.out.log
stdout_logfile=/var/log/supervisor/ordx-testnet4.out.log