# help/supervisor/conf/ordx.ini
[program:bitcoind-25.0-main]
priority=1
command=bitcoind -rpcthreads=128 -rpcworkqueue=512 -chain=main -conf=/data2/bitcoin-data/bitcoin.conf -datadir=/data2/bitcoin-data/main
autostart=false
autorestart=true
stopsignal=INT

[program:bitcoind-25.0-testnet]
priority=1
command=bitcoind -rpcthreads=128 -rpcworkqueue=512 -chain=test -conf=/data2/bitcoin-data/bitcoin.conf -datadir=/data2/bitcoin-data/
autostart=false
autorestart=true
stopsignal=INT

[program:ord0.16.0-main]
priority=2
command=ord0.16.0 --chain mainnet --bitcoin-rpc-url 192.168.1.103:8332 --data-dir /data2/ord-data/0.16.0/main --bitcoin-rpc-username jacky --bitcoin-rpc-password _RZekaGRgKQJSIOYi6vq0_CkJtjoCootamy81J2cDn0 --first-inscription-height 767430 --index-sats --minify server --http --http-port 82
autostart=false
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ord0.16.0-main.out.log
stdout_logfile=/var/log/supervisor/ord0.16.0-main.out.log

[program:ord0.16.0-testnet]
priority=2
command=ord0.16.0 --chain testnet --bitcoin-rpc-url 192.168.1.103:18332 --data-dir /data2/ord-data/0.16.0 --bitcoin-rpc-username jacky --bitcoin-rpc-password 123456 --first-inscription-height 2413343 --index-sats --minify server --http --http-port 80
autostart=false
autorestart=true
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ord0.16.0-testnet.out.log
stdout_logfile=/var/log/supervisor/ord0.16.0-testnet.out.log

[program:ordx-server-master-main]
priority=3
command=ordx-server -env /etc/ordx/dev-mainnet.env
autostart=false
autorestart=false
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ordx-server-main.out.log
stdout_logfile=/var/log/supervisor/ordx-server-main.out.log

[program:ordx-server-master-testnet]
priority=3
command=ordx-server -env /etc/ordx/dev-testnet.env
autostart=false
autorestart=false
stopsignal=INT
user=root
stderr_logfile=/var/log/supervisor/ordx-server-testnet.out.log
stdout_logfile=/var/log/supervisor/ordx-server-testnet.out.log